# docker/izus_model/supervisord.conf
; Supervisord config for the izus_model container.
;
; Runs the YOLO detector and the aic_model policy as two independent processes
; so each gets its own Python interpreter (no shared GIL) and the YOLO node is
; not recreated on every lifecycle reconfigure. Both inherit the Zenoh env
; exported by /entrypoint.sh.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:yolo_detector]
; use_sim_time:=true is REQUIRED — the detector stamps its TF broadcasts and
; image-throttling decisions with self.get_clock().now(). The policy node runs
; on sim time and ages those TFs against its own sim clock; without sim_time
; here every detection looks decades old and gets rejected (TF_TIMEOUT). All
; work must happen on sim time, not wall time.
command=pixi run --as-is ros2 run izus_policy test_yolov12_perception --ros-args -p use_sim_time:=true
directory=/ws_aic/src/aic
autorestart=true
startsecs=3
priority=10
stopwaitsecs=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:aic_model]
command=pixi run --as-is ros2 run aic_model aic_model --ros-args -p policy:=izus_policy.run_hybrid -p use_sim_time:=true
directory=/ws_aic/src/aic
autorestart=true
startsecs=5
priority=20
stopwaitsecs=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
